# kinguard-certhandler.cron
#!/bin/bash
# update certificates
/usr/share/kinguard-certhandler/letsencrypt.sh -u
